# images/basebuild/etc/supervisor/conf.d/consul.conf
[program:consul]
command=consul agent -config-dir /etc/consul.d -join %(ENV_HOST_IP)s
stopsignal=INT
priority=1
