# images/kafka/kafka-0.8.1.1/etc/supervisor/conf.d/kafka.conf
[program:kafka]
command=/usr/local/kafka/kafka-start.sh
stopsignal=INT
priority=1
