# images/cassandra/cassandra-2.0.9/etc/supervisor/conf.d/cassandra.conf
[program:cassandra]
command=/usr/local/cassandra/bin/start.sh
stopsignal=INT
priority=1
