# images/consul/supervisord.conf
[supervisord]
nodaemon=true

[program:consul]
command=/bin/bash -c "/opt/consul/consul agent -config-dir=/data/consul/config"